# quadlet/lifetrack-postgres.container
[Unit]
Description=LifeTrack PostgreSQL Database
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=lifetrack_postgres
Image=docker.io/library/postgres:16-alpine
Environment=POSTGRES_DB=${POSTGRES_DB}
Environment=POSTGRES_USER=${POSTGRES_USER}
Environment=POSTGRES_PASSWORD=${POSTGRES_PASSWORD}
Volume=lifetrack-postgres-data.volume:/var/lib/postgresql/data
Network=lifetrack.network
PublishPort=5435:5432
HealthCmd=pg_isready
HealthInterval=10s
HealthTimeout=5s
HealthRetries=5

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target
