# quadlet/lifetrack-frontend.container
[Unit]
Description=LifeTrack Frontend Web Application
After=network-online.target lifetrack-backend.service
Wants=network-online.target
BindsTo=lifetrack-backend.service

[Container]
ContainerName=lifetrack_frontend
Image=localhost/lifetrack-frontend:latest
Environment=NEXT_PUBLIC_API_URL=/query
Environment=BACKEND_URL=http://lifetrack_backend:8080
Environment=NODE_ENV=production
Network=lifetrack.network
HealthCmd=wget --no-verbose --tries=1 --spider http://localhost:3000
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Restart=unless-stopped
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target
