# quadlet/lifetrack-bot.container
[Unit]
Description=LifeTrack Telegram Bot with LLM
After=network-online.target lifetrack-backend.service
Wants=network-online.target
BindsTo=lifetrack-backend.service

[Container]
ContainerName=lifetrack_bot
Image=localhost/lifetrack-bot:latest
Environment=TELEGRAM_BOT_TOKEN=${TELEGRAM_BOT_TOKEN}
Environment=BACKEND_URL=http://lifetrack_backend:8080/query
Environment=SERVICE_JWT=${SERVICE_JWT}
Environment=FILE_STORAGE_PATH=/data/files
Environment=FILE_TEMP_PATH=/tmp/lifetrack_files
Volume=%h/lifetrack/models:/app/models:ro
Volume=%h/lifetrack/backend/graph:/backend/graph:ro
Volume=%h/lifetrack/bot/data:/app/data
Volume=%h/lifetrack/data/files:/data/files
Network=lifetrack.network

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target
