# quadlet/lifetrack-backend.container
[Unit]
Description=LifeTrack Backend GraphQL API
After=network-online.target lifetrack-postgres.service
Wants=network-online.target
BindsTo=lifetrack-postgres.service

[Container]
ContainerName=lifetrack_backend
Image=localhost/lifetrack-backend:latest
Environment=DATABASE_URL=postgres://${POSTGRES_USER}:${POSTGRES_PASSWORD}@lifetrack_postgres:5432/${POSTGRES_DB}?sslmode=disable
Environment=JWT_SECRET=${JWT_SECRET}
Environment=PORT=8080
Environment=ENV=production
Environment=FILE_STORAGE_PATH=/data/files
Volume=%h/lifetrack/data/files:/data/files
Network=lifetrack.network
HealthCmd=wget --no-verbose --tries=1 --spider http://localhost:8080/health
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target
