# quadlet/lifetrack-caddy.container
[Unit]
Description=LifeTrack Caddy Reverse Proxy
After=network-online.target lifetrack-backend.service lifetrack-frontend.service
Wants=network-online.target
BindsTo=lifetrack-backend.service
BindsTo=lifetrack-frontend.service

[Container]
ContainerName=lifetrack_caddy
Image=docker.io/library/caddy:2-alpine
Environment=ALLOWED_ORIGINS=${ALLOWED_ORIGINS}
Volume=%h/lifetrack/Caddyfile:/etc/caddy/Caddyfile:ro
Volume=lifetrack-caddy-data.volume:/data
Volume=lifetrack-caddy-config.volume:/config
Volume=lifetrack-caddy-logs.volume:/var/log/caddy
Network=lifetrack.network
PublishPort=8001:8001
HealthCmd=wget --no-verbose --tries=1 --spider http://localhost:8001/health
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target
